# ansible/roles/signal/templates/signal-cli.service.j2
[Unit]
Description=signal-cli JSON-RPC daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=runner
ExecStart=/usr/local/bin/signal-cli -a {{ nemoclaw_signal_account }} daemon --http --receive-mode=on-connection
WorkingDirectory=/var/lib/signal-cli
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
